# autostart/latte-dock.sh.desktop
[Desktop Entry]
Exec=/home/mod3sty/Scripts/latte-dock.sh
Icon=dialog-scripts
Name=latte-dock.sh
Path=
Type=Application
X-KDE-AutostartScript=true
